# BUILD/supervisord.conf
[supervisord]
nodaemon=true

[program:mageni-vscand]
command=/usr/local/sbin/mageni-scanner-start

[program:mageni-db]
command=/usr/local/sbin/mageni-sqlite-start

[program:mageni-web]
command=/usr/sbin/nginx -g 'daemon on; master_process on;'

[program:mageni-php]
command=service php8.0-fpm start
